# deploy/visibility-dashboard.service
[Unit]
Description=LLM Visibility Dashboard (FastAPI / uvicorn)
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/Uplers-llm-visibility/dashboard
EnvironmentFile=/home/ubuntu/Uplers-llm-visibility/.env
ExecStart=/home/ubuntu/Uplers-llm-visibility/dashboard/venv/bin/uvicorn app:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=10
StandardOutput=append:/var/log/visibility-dashboard.log
StandardError=append:/var/log/visibility-dashboard.log

[Install]
WantedBy=multi-user.target
